# packaging/linux/abacus.service
[Unit]
Description=Abacus
After=network.target
Wants=network.target

[Service]
Type=exec
ExecStart=/usr/local/bin/abacus
DynamicUser=true
StateDirectory=abacus
WorkingDirectory=%S/abacus
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
UMask=0077
StandardOutput=journal
StandardError=journal
SyslogIdentifier=abacus

[Install]
WantedBy=multi-user.target
